# deploy/house_shop.ini
[program:house_shop]
directory = /home/bezrukov/Envs/house_shop/source
user = bezrukov
command = /home/bezrukov/Envs/house_shop/source/gunicorn.sh
stdout_logfile = /home/bezrukov/Envs/house_shop/logs/supervisor_stdout.log
stderr_logfile = /home/bezrukov/Envs/house_shop/logs/supervisor_stderr.log
